# etc/supervisord.conf.d/suricata.ini
[program:suricata]
command=/usr/sbin/suricata %(ENV_SURICATA_ARGS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
